# claude-automations/launchd/com.2lines.qmd-indexer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.2lines.qmd-indexer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/john/Documents/Workspace/2Lines/knowledge-base &amp;&amp; /Users/john/.cargo/bin/qmd update &amp;&amp; /Users/john/.cargo/bin/qmd embed</string>
    </array>

    <!-- Run at 3:30 AM every day (before morning workflow) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/john/.cargo/bin:/Users/john/.local/bin</string>
        <key>HOME</key>
        <string>/Users/john</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/john/logs/claude-automations/launchd-qmd-indexer.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/john/logs/claude-automations/launchd-qmd-indexer.log</string>
</dict>
</plist>
