# claude-automations/launchd/com.2lines.granola-reminder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.2lines.granola-reminder</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/osascript</string>
        <string>-e</string>
        <string>display notification "Ensure Granola is connected before 8pm timesheet sync" with title "🥣 Granola Auth Check" sound name "Blow"</string>
    </array>

    <!-- Run at 7:00 PM every day (1 hour before evening workflow) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
</dict>
</plist>
